# crates/stat-tracker/assets/scuffed-stat-tracker.service
[Unit]
Description=Scuffed Stat Tracker — OW2 personal stats capture daemon
After=graphical-session.target
PartOf=graphical-session.target
# A daemon that fails to start on every attempt (e.g. a broken config or a
# genuinely-stuck lock) must land in `failed` and stay visible — NOT loop on
# Restart=on-failure forever. Give up after 5 failures inside 5 minutes so the
# GUI/journal show a clear stopped state instead of a silent 300+ restart churn.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/bin/scuffed-stat-tracker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# Yield CPU and disk to the game: OCR bursts run niced and in a
# low-weight cgroup so Overwatch keeps frame-pacing priority (H2).
Nice=10
CPUWeight=30
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=graphical-session.target
